# systemd/screensaver-daemon.service
[Unit]
Description=Terminal Screensaver Daemon
Documentation=https://github.com/reisset/myscreensavers
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/myscreensavers/bin/screensaver-daemon
Restart=on-failure
RestartSec=5

# Optional: Customize idle timeout (in milliseconds)
# Default is 300000 (5 minutes)
Environment=IDLE_TIMEOUT_MS=300000
Environment=CHECK_INTERVAL=5

[Install]
WantedBy=graphical-session.target
